# dockerfiles/ingest/crontab
# Ingest every fifteen minutes
0 */5 * * * * * /app/parse_and_ingest.sh
