# dockerfiles/ingest/crontab
# Ingest every fifteen minutes
0 */3 * * * * * /app/parse_and_ingest.sh
