# dockerfiles/backup/crontab
# backup every hour
@hourly /app/backup.sh
